# linux/wayland/eim-text-expansion.service
[Unit]
Description=EIM Text Expansion Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/python3 /path/to/EIM_autokey_dotool_daemon.py --background
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
